# gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn Flask/app:run -b localhost:5000
directory=/deploy/app